# packages/rbus/rbus.service
[Unit]
Description=RBus IPC framework router daemon
After=network.target

[Service]
Type=forking
TimeoutSec=300
ExecStart=/usr/bin/rtrouted
ExecStopPost=/bin/sh -c 'rm -f /tmp/rtrouted*'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
