# packages/rbus-elements/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Reload systemd daemon
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null 2>&1 || true
            
            # Enable the service
            if ! systemctl is-enabled rbus-elements.service >/dev/null 2>&1; then
                systemctl enable rbus-elements.service >/dev/null 2>&1 || true
            fi
            
            # Start the service if rbus is running
            if systemctl is-active rbus.service >/dev/null 2>&1; then
                if ! systemctl is-active rbus-elements.service >/dev/null 2>&1; then
                    systemctl start rbus-elements.service >/dev/null 2>&1 || true
                fi
            fi
        fi
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
