# packages/parodus/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Reload systemd daemon
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null 2>&1 || true
            
            # Enable and start the service
            if ! systemctl is-enabled parodus.service >/dev/null 2>&1; then
                systemctl enable parodus.service >/dev/null 2>&1 || true
            fi
            
            # Start the service if not running
            if ! systemctl is-active parodus.service >/dev/null 2>&1; then
                systemctl start parodus.service >/dev/null 2>&1 || true
            fi
        fi
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
